# agent/systemd/jjudge.slice
[Unit]
Description=Slice for jjudge's sandboxes
